# crontab
0 2 * * * cd /app && python -m scripts.generate_report >> /var/log/cron.log 2>&1